# src/dionaea-docker/0.5/supervisor/dionaea.conf
[program:dionaea]
command=/opt/dionaea/bin/dionaea -u dionaea -g dionaea -c /opt/dionaea/etc/dionaea/dionaea.cfg
autorestart=true
