# old/src/dionaea-docker/0.5/supervisor/p0f.conf
[program:p0f]
command=p0f -i any -u dionaea -Q /tmp/p0f.sock -q -l
autorestart=true
